# file_daemon/install_weka.cron
TOKEN=<TOKEN>
* * * * * /path/to/install.sh
