# file_daemon/weka_testing.service
[Unit]
Description=Dummy Service

[Service]
ExecStart=/bin/sleep infinity

[Install]
WantedBy=multi-user.target